# battery.service
[Service]
ExecStart=battery_notifier
Restart=always
Type=simple

[Install]
WantedBy=default.target
