# battery.service
[Service]
ExecStart=battery-notification-daemon
Restart=always
Type=simple

[Install]
WantedBy=default.target
